# root/etc/s6-overlay/s6-rc.d/svc-wol-proxy/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Starting WOL-proxy"

exec \
    cd /app s6-setuidgid abc python3 WOL-proxy.py
